# config/app/gunicorn.conf
[program:gunicorn]
directory=/var/code
command=gunicorn -b :8888 python_graduate.wsgi:application
autostart=true
autorestart=true
stdout_logfile=/var/logs/gunicorn.log
redirect_stderr=true